# polymarket-focus-bridge.service
[Unit]
Description=Polymarket Focus Command Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/polymarket
Environment="PATH=/root/polymarket/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/polymarket/.venv/bin/python -m uvicorn focus_command_bridge:APP --host 127.0.0.1 --port 8012
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
